# checkbattery/crontab
DISPLAY=:0.0
XAUTHORITY=/home/matrix/.Xauthority

#runs check battery script each 5 minutes
*/5 * * * * /bin/sh /home/[USER]/check_battery.sh >> /home/[USER]/check_battery.log 2>&1

